# usr/lib/systemd/system/rkm_playMusic.service
[Unit]
Description=Tocar musica estipulada na saida analogica

[Service]
ExecStart=/usr/bin/rkm_playMusic

[Install]
WantedBy=multi-user.target
