# usr/lib/systemd/system/rkm_autoStart.service
[Unit]
Description="Script de AutoStart"

[Service]
ExecStart=/bin/rkm_autoStart

[Install]
WantedBy=multi-user.target
