# usr/lib/systemd/system/rkm_internet.service
[Unit]
Description=Verificar a conexao com a internet

[Service]
ExecStart=/usr/bin/rkm_internet

[Install]
WantedBy=multi-user.target
