# usr/lib/systemd/system/rkm_conexoesWifi.service
[Unit]
Description=monitorar as redes wifi disponiveis

[Service]
ExecStart=/usr/bin/rkm_conexoesWifi

[Install]
WantedBy=multi-user.target
